# wwwroot/infra/systemd/lowland-api.service
[Unit]
Description=Lowland API (Hulpveren backend)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/apps/backend/server
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /opt/apps/backend/server/index.js
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
